# Spectre/airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
# Spectre OS: No auto-login. User must authenticate.
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f \u' --noclear %I $TERM
