# airootfs/etc/profile.d/spectre-motd.sh
# Spectre OS - MOTD on login
[[ -f /etc/spectre/motd.sh ]] && /etc/spectre/motd.sh
