# airootfs/etc/systemd/system/spectre-firstboot.service
[Unit]
Description=Spectre OS - First boot setup
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spectre-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
